# scripts/pgBackup.timer
[Unit]
Description=Backup postgres timer
Requires=pgBackup.service

[Timer]
Unit=pgBackup.service
OnCalendar=20:00

[Install]
WantedBy=timers.target
